# vagrant/provision/crontab
10 * * * * /usr/bin/php /var/www/yii rating
